# templates/etc/systemd/system/activemq.service.erb
[Unit]
Description=Activemq Broker Service
After=network.target
Wants=activemq.service

[Service]
PrivateTmp=true
ExecStart=/opt/activemq/bin/activemq start
ExecReload=/opt/activemq/bin/activemq restart
ExecStop=/opt/activemq/bin/activemq stop
User=activemq
Group=activemq
Type=forking
Environment=JAVA_HOME=/usr/bin/java
PIDFile=/opt/activemq/data/activemq-<%= @hostname %>.pid
KillMode=process

[Install]
WantedBy=multi-user.target
