# templates/activemq.service.erb
[Unit]
Description=Apache ActiveMQ

[Service]
User=activemq
Type=forking
PIDFile=/opt/activemq/data/activemq.pid
<%- if @jmx_enabled -%>
Environment='ACTIVEMQ_SUNJMX_START=<%= @jmx_agent_opts %>' 'ACTIVEMQ_USER=activemq'
<%- end -%>
ExecStart=/opt/activemq/bin/activemq start --extdir /opt/activemq/lib/tipaas
ExecReload=/opt/activemq/bin/activemq restart --extdir /opt/activemq/lib/tipaas
ExecStop=/opt/activemq/bin/activemq kill
LimitNPROC=64000
LimitNOFILE=64000

[Install]
Alias=activemq
