# templates/activemq.service.erb
[Unit]
Description=Apache ActiveMQ

[Service]
Type=forking
PIDFile=/opt/activemq/data/activemq.pid
ExecStart=/opt/activemq/bin/activemq start --extdir lib/tipaas
ExecReload=/opt/activemq/bin/activemq restart --extdir lib/tipaas
ExecStop=/opt/activemq/bin/activemq kill
LimitNPROC=64000
LimitNOFILE=64000

[Install]
Alias=activemq
